# deploy/freelo_scrape.conf
description "freelo server"
start on runlevel []
stop on runlevel [06]

exec /usr/bin/python /home/ubuntu/lol-master-ease/data_mining/scraper.py
